# python-fitbit/crontab
*/1 * * * * /Users/anthony/Desktop/python_api_puller/python-fitbit/token_alive.py
